# projects/playground/systemd/fix-overlay-order.service.in
[Unit]
Description=Fixes bug in initramfs
After=syslog.target persistent.mount

[Service]
ExecStart=@C15_INSTALL_PATH@/scripts/fixOverlayOrder.sh

[Install]
WantedBy=multi-user.target
