# projects/bbb/bridge/systemd/install-update-from-epc.service.in
[Unit]
Description=Copy Update From EPC and Execute

[Service]
Type=simple
ExecStart=@C15_INSTALL_PATH@/scripts/install-update-from-epc.sh