# projects/bbb/bridge/systemd/midi-over-tcp.service.in
[Unit]
Description=Share external MIDI device over network
After=aseqnet.service

[Service]
# FIND PORTS FOR MIDI DEVICES
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment EXTERNAL_MIDI_DEVICE=$(aplaymidi -l |  grep -v 'Net Client' | grep -o '[0-9]*:[0-9]*' | head -n1)"
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment NETWORK_MIDI_DEVICE=$(aplaymidi -l |  grep 'Net Client' | grep -o '[0-9]*:[0-9]*' | head -n1)"

# SETUP EXTERNAL MIDI CONNECTION VIA NETWORK
ExecStartPre=/usr/bin/aconnect -x
ExecStart=/usr/bin/aconnect ${EXTERNAL_MIDI_DEVICE} ${NETWORK_MIDI_DEVICE}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
