# projects/bbb/bridge/systemd/redirect-websocket-to-recorder.service.in
[Unit]
Description=Redirect recorder websocket traffic to ePC
After=syslog.target network.target

[Service]
ExecStart=/bin/socat TCP-LISTEN:@RECORDER_WEBSOCKET_PORT@,fork TCP:192.168.10.10:@RECORDER_WEBSOCKET_PORT@
Restart=always
RestartSec=1
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target 