# projects/bbb/bridge/systemd/redirect-http-to-playground.service.in
[Unit]
Description=Redirect playground http traffic to ePC
After=syslog.target network.target

[Service]
ExecStart=/bin/socat TCP-LISTEN:@PLAYGROUND_HTTPSERVER_PORT@,fork TCP:192.168.10.10:@PLAYGROUND_HTTPSERVER_PORT@
Restart=always
RestartSec=1
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target