# projects/bbb/bridge/systemd/reset-playcontroller.service.in
[Unit]
Description=safety reset the playcontroller (primarily needed for older motherboards)
After=systemd-modules-load.service
Before=install-update-from-usb.service

[Service]
Type=oneshot
ExecStart=-@C15_INSTALL_PATH@/scripts/reset-playcontroller.sh
