# projects/bbb/bridge/systemd/reset-playcontroller.service.in
[Unit]
Description=safety reset the playcontroller, for older motherboards
After=systemd-modules-load.service
Before=install-update-from-usb.service

[Service]
Type=oneshot
ExecStart=/usr/bin/playcontroller reset 5
