# projects/bbb/midi-over-ip/systemd/midi-over-ip.service.in
[Unit]
Description=Nonlinear-Labs Midi Bridge

[Service]
ExecStart=@C15_INSTALL_PATH@/midi-over-ip/midi-over-ip -a 192.168.10.10

[Install]
WantedBy=multi-user.target
