# projects/bbb/midi-over-ip/systemd/midi-over-ip.service.in
[Unit]
Description=Nonlinear-Labs Midi Bridge

[Service]
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment MIDI_DEVICE=$(amidi -l | tail -n 1 | grep -o 'hw:[0-9,]*')"
ExecStart=@C15_INSTALL_PATH@/midi-over-ip/midi-over-ip -m ${MIDI_DEVICE} -a 192.168.10.10

[Install]
WantedBy=multi-user.target
