# projects/epc/audio-engine/systemd/audio-engine.service.in
[Unit]
Description=Nonlinear-Labs C15 Audio Engine
After=systemd-modules-load.service sound.target fixCpuBindings.service external-midi.service
StartLimitIntervalSec=0
StartLimitInterval=0

[Service]
WorkingDirectory=@C15_INSTALL_PATH@/audio-engine

# WAIT FOR PLAYCONTROLLER
ExecStartPre=/usr/bin/bash -c "while [ ! -f /proc/asound/EMPHASE/id ]; do sleep 0.1; done"

# WAIT FOR AUDIO DEVICE
ExecStartPre=/usr/bin/bash -c "while [ ! -f /proc/asound/PCH/id ]; do sleep 0.1; done"

# FULL VOLUME
ExecStartPre=amixer -D hw:PCH set Master 100%

# SETUP FRAME SIZE DEPENDENT ON CPU
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment FRAMES_PER_PERIOD=49"
ExecStartPre=/usr/bin/bash -c "if cat /proc/cpuinfo | grep \"model name\" | grep \"7010\" > /dev/null; then /usr/bin/systemctl set-environment FRAMES_PER_PERIOD=49; fi"
ExecStartPre=/usr/bin/bash -c "if cat /proc/cpuinfo | grep \"model name\" | grep \"5010\" > /dev/null; then /usr/bin/systemctl set-environment FRAMES_PER_PERIOD=64; fi"

# FIND PORTS AND NAMES FOR EXTERNAL MIDI DEVICES
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment VIRTUAL_MIDI_DEVICE=$(aplaymidi -l | grep 'Virtual' | head -n1 | grep -o '[0-9]*:[0-9]*')"
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment NETWORK_MIDI_DEVICE=$(aplaymidi -l | grep 'Network' | head -n1 | grep -o '[0-9]*:[0-9]*')"
ExecStartPre=/usr/bin/bash -c "/usr/bin/systemctl set-environment VIRTUAL_MIDI_DEVICENAME=$(amidi -l | grep 'Virtual' | head -n1 | grep -o 'hw:[^ ]*')"

# SETUP EXTERNAL MIDI CONNECTION VIA NETWORK
ExecStartPre=/usr/bin/aconnect -x
ExecStartPre=/usr/bin/aconnect ${NETWORK_MIDI_DEVICE} ${VIRTUAL_MIDI_DEVICE}

# START THE ENGINE
ExecStart=@C15_INSTALL_PATH@/audio-engine/audio-engine -t hw:EMPHASE -m ${VIRTUAL_MIDI_DEVICENAME} -a hw:PCH -n 3 -s ${FRAMES_PER_PERIOD} -r 48000 --playground-host=localhost --heartbeat=hw:EMPHASE

Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
