# playground/systemd/playground.service.in
[Unit]
Description=Nonlinear-Labs Playground
After=syslog.target network.target systemd-modules-load.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/nonlinear/playground/playground --bbbb=192.168.8.2
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
