# projects/epc/audio-engine/system/audio-engine.service.in
[Unit]
Description=Nonlinear-Labs C15 Audio Engine
After=systemd-modules-load.service sound.target alsa-restore.service
StartLimitIntervalSec=0
StartLimitInterval=0

[Service]
WorkingDirectory=@C15_INSTALL_PATH@/audio-engine

# WAIT FOR PLAYCONTROLLER
ExecStartPre=/usr/bin/bash -c "while [ ! -f /proc/asound/EMPHASE/id ]; do sleep 0.1; done"

# WAIT FOR AUDIO DEVICE
ExecStartPre=/usr/bin/bash -c "while [ ! -f /proc/asound/PCH/id ]; do sleep 0.1; done"

# START THE ENGINE
ExecStart=@C15_INSTALL_PATH@/audio-engine/audio-engine -t hw:EMPHASE -a hw:PCH -n 8 -s 16 -r 48000 --playground-host=localhost --heartbeat=hw:EMPHASE

Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
