# projects/bbbb/systemd/chkphy.service.in
[Unit]
Description=Check BBB phy on boot and repower if PHY failed
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=@C15_INSTALL_PATH@/scripts/bbb-check-phy.sh
#RemainAfterExit=yes
TimeoutSec=0
