# projects/bbbb/systemd/gettimefromepc.service.in
[Uint]
Description=ePC Time via sshpass
After=syslog.target network.target systemd-modules-load.service

[Service]
ExecStart=/usr/C15/scripts/gettimefromepc.sh
Restart=always
RestartSec=20
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
