# projects/bbbb/systemd/handle-mmc-device-appeared-twice.service.in
[Unit]
Description=Handle "MMC Device appeared twice"-bug
After=bbbb.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@C15_INSTALL_PATH@/scripts/handle-mmc-device-appeared-twice.sh
