# projects/bbbb/systemd/install-update-from-usb.service.in
[Unit]
Description=Copy Update From USB Stick
After=syslog.target network.target chkphy.service systemd-modules-load.service

[Service]
ExecStart=@C15_INSTALL_PATH@/scripts/install-update.sh

[Install]
WantedBy=multi-user.target

