# projects/bbbb/systemd/reset-lpc.service.in
[Unit]
Description=safety reset the LPC, for older motherboards
After=systemd-modules-load.service
Before=install-update-from-usb.service

[Service]
Type=oneshot
ExecStart=/usr/bin/lpc reset 5
